# src/99-fit-sync.rules.in
# udev -*- mode: conf -*- file for my local linux (Raspberry PI) installation.

# Keep all the things in one place, even though they are not really related.
# try "man udev" for documentation

# HINT: look for the device in /var/log/syslog and run  "sudo udevadm test
# $DEVNAME" to test these rules.

### ANT Usb stick rules

# Change the mode of the Dynastream ANT UsbStick2 so all users can read and
# write to it.

# If you do not want all users to be able to read and write to your usb stick
# you might want to adjust the mode to 0660 and restrict it to a group of
# users with GROUP="ant" or the mode to 0600 and restrict it to a single user
# with OWNER="user".

SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="0fcf", ATTR{idProduct}=="1008", MODE="0666", SYMLINK+="ttyANT2"
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="0fcf", ATTR{idProduct}=="1009", MODE="0666", SYMLINK+="ttyANT3"

SUBSYSTEM=="usb", ACTION=="add", ATTR{interface}=="MTP" RUN+="@@BINDIR@@/on-mtp-added $devpath"
SUBSYSTEM=="usb", ACTION=="remove", RUN+="@@BINDIR@@/on-mtp-removed $kernel"

### Mount a Garmin and Fly6 camera (as a disk) and do the right thing for each
### of them.

# we try to mount devices appear as a SCSI disks (sd*) with or without numbers.
# Garmin has the file system on /dev/sda, Fly6 on /dev/sda1

KERNEL=="sd[a-z]", GOTO="start_automount"
KERNEL=="sd[a-z][0-9]", GOTO="start_automount"
GOTO="end_automount"

LABEL="start_automount"
IMPORT{program}="/sbin/blkid -o udev -p /dev/%k"
ENV{ID_FS_LABEL}=="GARMIN", GOTO="garmin_automount"
GOTO="end_automount"

### Mount a Garmin device and start synching files

LABEL="garmin_automount"

ACTION=="add", RUN+="/bin/systemctl start mount-garmin@%k.service"
ACTION=="add", RUN+="/bin/systemctl start fit-sync-usb.service"
# 25/06/11 Don't run this, it broke somwhere along the way and I don't need it
# anyway

# ACTION=="add", RUN+="/bin/systemctl start fit-sync-epo.service"

# Run umount when the device is removed (this will only useful if the device
# is prematurely removed.)
ACTION=="remove", RUN+="/bin/systemctl stop mount-garmin@%k.service"

GOTO="end_automount"

LABEL="end_automount"
