# src/fit-sync-usb.service.in
; -*- mode: conf -*-

[Unit]
Description=Download FIT files from Garmin devices (USB version)

[Service]
Type=forking
PIDFile=/run/fit-sync/fit-sync-usb.pid
ExecStart=@@PREFIX@@/bin/fit-sync-usb -d /media/garmin/GARMIN
Restart=no
WorkingDirectory=/media/garmin
User=pi
Group=users
