# src/fit-sync-epo.service.in
; -*- mode: conf -*-

[Unit]
Description=Update EPO data for Garmin devices

[Service]
Type=simple
ExecStart=@@BINDIR@@/fit-sync-epo -s REMOTESW/EPO.BIN
Restart=no
WorkingDirectory=@@GARMIN_MOUNT_POINT@@/GARMIN
User=@@USER@@
Group=@@GROUP@@
