# src/fit-sync-usb.service.in
; -*- mode: conf -*-

[Unit]
Description=Download FIT files from Garmin devices (USB version)
After=fit-sync-setup.service

[Service]
Type=forking
PIDFile=/run/fit-sync/fit-sync-usb.pid
ExecStart=@@BINDIR@@/fit-sync-usb -d /media/garmin/GARMIN
Restart=no
WorkingDirectory=/media/garmin
User=@@USER@@
Group=@@GROUP@@
