# src/fit-sync-epo.service.in
; -*- mode: conf -*-

[Unit]
Description=Update EPO data for Garmin devices

[Service]
Type=simple
ExecStart=@@BINDIR@@/fit-sync-epo -s /media/garmin/GARMIN/REMOTESW/EPO.BIN
Restart=no
WorkingDirectory=/media/garmin
User=@@USER@@
Group=@@GROUP@@
