# src/mount-fr965.service.in
; -*- mode: conf -*-
[Unit]
Description=Mount the Forerunner 965 MTP device

[Service]
Type=forking
ExecStart=@@BINDIR@@/mount-fr965
ExecStop=/usr/bin/fusermount -u /media/fr965
Restart=no
WorkingDirectory=/
User=root
Group=root
