# src/sync-fr965.service.in
; -*- mode: conf -*-

[Unit]
Description=Download FIT files from a FR965
After=fit-sync-setup.service
Requires=mount-fr965.service
After=mount-fr965.service

[Service]
Type=forking
PIDFile=/run/fit-sync/fit-sync-fr965.pid
ExecStart=@@BINDIR@@/fit-sync-usb -p /run/fit-sync/fit-sync-fr965.pid -d .
Restart=no
WorkingDirectory=/media/fr965/Internal Storage/GARMIN
User=@@USER@@
Group=@@GROUP@@
