# src/sync-fr945.service.in
; -*- mode: conf -*-

[Unit]
Description=Download FIT files from a FR945
After=fit-sync-setup.service
Requires=mount-fr945.service
After=mount-fr945.service

[Service]
Type=forking
PIDFile=/run/fit-sync/fit-sync-fr945.pid
ExecStart=@@BINDIR@@/fit-sync-usb -p /run/fit-sync/fit-sync-fr945.pid -d .
Restart=no
WorkingDirectory=/media/fr945/Internal Storage/GARMIN
User=@@USER@@
Group=@@GROUP@@
